# debian/cron.d
#
# Regular cron jobs for the ltposter package
#
0 4	* * *	root	rm -rf /usr/share/ltposter/records/*
